# home/Library/LaunchAgents/osx.app.1password.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>osx.app.1password</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Applications/1Password 6.app/Contents/MacOS/1Password 6</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
