# home/Library/LaunchAgents/osx.app.fantastical.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>osx.app.fantastical</string>
    <key>LaunchOnlyOnce</key>
    <true/>
    <key>ProgramArguments</key>
    <array>
      <string>sh</string>
      <string>-c</string>
      <string>sleep 10 &amp;&amp; /Applications/Fantastical\ 2.app/Contents/MacOS/Fantastical\ 2</string>
    </array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
