# home/Library/LaunchAgents/osx.app.drop-zone.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>osx.app.drop-zone</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-c</string>
		<string>BASH_NO=no . /etc/profile
        /Applications/Dropzone\ 3.app/Contents/MacOS/Dropzone\ 3</string>
	</array>
	<key>RunAtLoad</key>
    <true/>
</dict>
</plist>
