# home/Library/LaunchAgents/osx.app.control-plane.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>osx.app.control-plane</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-c</string>
		<string>BASH_NO=no . /etc/profile &amp;&amp; sleep 6 &amp;&amp; /Applications/ControlPlane.app/Contents/MacOS/ControlPlane</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
