# home/Library/LaunchAgents/osx.app.menu-pop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<true/>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>osx.app.menu-pop</string>
	<key>ProgramArguments</key>
	<array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>if [ -d "/Applications/MenuPop.app" ]; then /Applications/MenuPop.app/Contents/MacOS/MenuPop; fi</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
