# home/Library/LaunchAgents/osx.environment.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>osx.environment</string>
  <key>ProgramArguments</key>
  <!-- PATH does not work, bug at apple side -->
  <array>
    <string>sh</string>
    <string>-c</string>
    <string>
        /bin/launchctl setenv BREW_HOME /usr/local
        /bin/launchctl setenv GOPATH /usr/local/go
        /bin/launchctl setenv NVM_DIR /usr/local/nvm
        /bin/launchctl setenv NVM_TARGET /usr/local/opt/nvm
      </string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardErrorPath</key>
  <string>/var/log/osx.environment.stderr</string>
  <key>StandardOutPath</key>
  <string>/var/log/osx.environment.stdout</string>
</dict>
</plist>
