# home/Library/LaunchAgents/osx.app.ghost-note.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<false/>
	<key>Label</key>
	<string>osx.app.ghost-note</string>
	<key>ProgramArguments</key>
	<array>
		<string>sh</string>
		<string>-c</string>
		<string>sleep 10 &amp;&amp; /Applications/GhostNote.app/Contents/MacOS/GhostNote</string>
	</array>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
